# src/.config/systemd/user/mbsync.timer
[Unit]
Description=Sync Mail Schedule

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
Unit=mbsync.service

[Install]
WantedBy=timers.target
